# systemd/user/udiskie.service
[Unit]
Description=Udiskie auto mount utility

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
